# roles/paleo_static/templates/paleo-static.service.j2
# copyright Utrecht University
# {{ ansible_managed }}
[Unit]
Description=Paleo Earth static site
# Docker service requires network up
After=network.target docker.service
# If docker dies, paleo-static dies
Requires=docker.service

[Service]
User={{ paleo_user }}
Group={{ paleo_user }}
WorkingDirectory=/home/{{ paleo_user }}/paleo-ansible/docker/static
# --no-recreate: don't bounce running containers on restart. MVP-only; drop it
# in prod (pinned images make `up` recreate only on real changes).
ExecStart=/usr/bin/docker compose up --no-color --no-recreate
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0
Restart=always

[Install]
WantedBy=multi-user.target
